# supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock

[supervisord]
nodaemon=false
logfile=/var/log/supervisor/supervisord.log
logfile_maxbytes=5MB
loglevel=info

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock ; use a unix:// URL  for a unix socket

[program:thumbor]
command=thumbor --ip=127.0.0.1 --port=800%(process_num)s --conf=/etc/thumbor.conf --log-level=debug
process_name=thumbor800%(process_num)s  
numprocs=4  
autostart=true  
autorestart=true  
startretries=3  
stdout_events_enabled=true
stderr_events_enabled=true

[program:nginx]
command=/usr/sbin/nginx
stdout_events_enabled=true
stderr_events_enabled=true

[eventlistener:stdout]
command = supervisor_stdout
buffer_size = 100
events = PROCESS_LOG
result_handler = supervisor_stdout:event_handler
